# crontab/prueba.sh
#!/bin/bash
echo "hola mundo" >> /code/log