# crontab/root
*/1 * * * * /bin/sh /code/ubicar_conductores.sh
*/1 * * * * /bin/sh /code/prueba.sh
