# crontab/root
*/1 * * * * /code/ubicar_conductores.sh
*/1 * * * * /code/prueba.sh
